# contrib/openclaw/openclaw-nextcloud-talk-poller.service
[Unit]
Description=OpenClaw Nextcloud Talk command fallback poller
Documentation=https://github.com/ShapeShifter499/agentcommands
After=network-online.target openclaw-gateway.service
Wants=network-online.target openclaw-gateway.service

[Service]
Type=simple
EnvironmentFile=%h/.config/agentcommands/nextcloud-talk-poller.env
ExecStart=%h/.local/bin/nextcloud-talk-poller --watch
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
